# ansible/files/fdo/etc/systemd/system/fdo-ansible-glue.service
[Unit]
Description=FDO Ansible Glue Service

[Service]
Type=oneshot
ExecStart=/var/local/glue-init.sh